# Sistema VPS/baseline_backups/deploy_systemd_bak/smc-opportunity-scanner.service
[Unit]
Description=Opportunity Scanner — Motor de scan de oportunidades (shadow-only)
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=bimaq
WorkingDirectory=/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0/venv/bin:/usr/bin:/bin
ExecStart=/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0/venv/bin/python3 tools/run_opportunity_scanner.py --loop --apply --heartbeat
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
