# Sistema VPS/baseline_backups/systemd_bak/smc-ngrok.service
[Unit]
Description=SMC ngrok Tunnel para Streamlit
After=smc-streamlit.service
Requires=smc-streamlit.service

[Service]
Type=simple
User=bimaq
ExecStartPre=/bin/sleep 3
ExecStart=/home/bimaq/bin/ngrok http 8501 --log-format=logfmt --log=/tmp/ngrok.log
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
