# Sistema VPS/baseline_backups/systemd_bak/smc-mt5linux-b3.service
[Unit]
Description=SMC mt5linux Bridge B3 (porta 11000)
After=smc-mt5-b3-terminal.service
Requires=smc-mt5-b3-terminal.service

[Service]
Type=simple
User=bimaq
Environment=DISPLAY=:99
Environment=WINEPREFIX=/home/bimaq/.wine-mt5-xp
Environment=WINEDEBUG=fixme-all
WorkingDirectory=/home/bimaq/.wine-mt5-xp/drive_c
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/wine python -m mt5linux --port 11000 --host 0.0.0.0
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
