# Sistema VPS/baseline_backups/systemd_bak/smc-forex-robot.service
[Unit]
Description=SMC Trader Forex Robot (XAUUSD, BTCUSD, EURUSD, USDJPY, ETHUSD, XAGUSD)
After=network.target

[Service]
Type=simple
User=bimaq
WorkingDirectory=/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0
ExecStart=/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0/venv/bin/python run_forex.py
Restart=always
RestartSec=10
StandardOutput=append:/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0/runtime/logs/forex_robot.log
StandardError=append:/home/bimaq/projetos/SMC_Trader_System_7_0/SMC_Trader_System 7.0/runtime/logs/forex_robot.log

[Install]
WantedBy=default.target
