# Sistema VPS/baseline_backups/systemd_bak/smc-mt5-b3-terminal.service
[Unit]
Description=SMC MT5 B3 Terminal (XP Investimentos)
After=smc-xvfb.service
Requires=smc-xvfb.service

[Service]
Type=simple
User=bimaq
Environment=DISPLAY=:99
Environment=WINEPREFIX=/home/bimaq/.wine-mt5-xp
Environment=WINEDEBUG=fixme-all
ExecStart=/usr/bin/wine "C:\\Program Files\\MetaTrader 5 XP\\terminal64.exe"
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
